# mac/kanata/.config/kanata/com.example.kanata.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.kanata</string>
  <!-- Use Homebrew symlink: survives brew upgrade so Input Monitoring permissions persist -->
  <!-- 5s delay gives Karabiner DriverKit time to initialize on boot -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>sleep 5 &amp;&amp; exec /opt/homebrew/bin/kanata -c /Users/cyperx/.config/kanata/kanata.kbd --port 10000</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>5</integer>
  <key>StandardOutPath</key>
  <string>/tmp/kanata.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/kanata.err.log</string>
</dict>
</plist>
